# ops/com.nyel.harness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.nyel.harness</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>caffeinate -is mix phx.server</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/nyelbangash/Documents/agent-harness/harness</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/nyelbangash/.harness/logs/harness.out.log</string>
  <key>StandardErrorPath</key><string>/Users/nyelbangash/.harness/logs/harness.err.log</string>
  <key>EnvironmentVariables</key>
  <dict><key>MIX_ENV</key><string>prod</string></dict>
</dict></plist>
